# gpu_stats.cron
*/1 * * * * root /bin/bash /usr/local/bin/gpu_stats.sh > /tmp/gpu_stats.log 2>&1
